# pulldocker.timer
[Unit]
Description=Pulldocker

[Timer]
OnCalendar=*-*-* *:15:00

[Install]
WantedBy=timers.target
